# instructlab-ui.fixed.container
# instructlab-ui.container (fixed)
# Purpose: JupyterLab UI using the same image; publishes localhost:8777 -> container:8888

[Unit]
Description=InstructLab UI (JupyterLab) - CUDA 12.6 cuDNN UBI9
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/instructlab/instructlab:latest

ContainerName=instructlab-ui


# Optional GPU access for notebooks
Environment=NVIDIA_VISIBLE_DEVICES=all
Environment=NVIDIA_DRIVER_CAPABILITIES=compute,utility
AddDevice=nvidia.com/gpu=all

# SELinux
SecurityLabelDisable=true
# Volume=%h/instructlab-data:/workspace:z
# Volume=%h/instructlab-data/.cache:/root/.cache:z

# Persistence
Volume=%h/instructlab-data:/workspace
Volume=%h/instructlab-data/.cache:/root/.cache

# Networking: bind only to loopback by default
PublishPort=127.0.0.1:8777:8888/tcp

WorkingDir=/workspace

# Jupyter config
Environment=JUPYTER_TOKEN=instructlab
Environment=JUPYTER_NOTEBOOK_DIR=/workspace

# Install minimal deps and start JupyterLab
Command=/bin/bash -lc "  dnf -y install python3.11 python3.11-pip git &&   python3.11 -m pip install --upgrade pip &&   python3.11 -m pip install jupyterlab instructlab &&   jupyter lab     --ServerApp.token=${JUPYTER_TOKEN}     --ServerApp.root_dir=${JUPYTER_NOTEBOOK_DIR}     --ServerApp.allow_remote_access=1     --ServerApp.ip=0.0.0.0     --ServerApp.port=8888     --ServerApp.open_browser=False "

[Service]
TimeoutStopSec=30
ExecStartPre=/usr/bin/bash -lc 'while true; do /usr/bin/nvidia-smi -L && exit 0; sleep 1; done'
Restart=always

[Install]
WantedBy=default.target
