# sdg-hub-ui.fixed.container
# sdg-hub-ui.container (fixed)
# Purpose: SDG Hub UI (JupyterLab). Shares the same persistent path.

[Unit]
Description=SDG Hub UI (JupyterLab) - CUDA 12.6 cuDNN UBI9
After=network-online.target
Wants=network-online.target

[Container]
# If SDG Hub needs CUDA base directly, keep this; otherwise align with instructlab image.
Image=docker.io/nvidia/cuda:12.6.0-cudnn-runtime-ubi9

ContainerName=sdg-hub-ui
Restart=always

# GPU hooks
Environment=NVIDIA_VISIBLE_DEVICES=all
Environment=NVIDIA_DRIVER_CAPABILITIES=compute,utility
AddDevice=nvidia.com/gpu=all

# SELinux
SecurityLabelDisable=true
# Volume=%h/instructlab-data:/workspace:z
# Volume=%h/instructlab-data/.cache:/root/.cache:z

# Persistence
Volume=%h/instructlab-data:/workspace
Volume=%h/instructlab-data/.cache:/root/.cache

# Networking
PublishPort=127.0.0.1:8999:8888/tcp

WorkingDir=/workspace

# UI token
Environment=JUPYTER_TOKEN=sdg
Environment=JUPYTER_NOTEBOOK_DIR=/workspace

# Launch JupyterLab
Command=/bin/bash -lc "  dnf -y install python3.11 python3.11-pip git &&   python3.11 -m pip install --upgrade pip &&   python3.11 -m pip install jupyterlab &&   jupyter lab     --ServerApp.token=${JUPYTER_TOKEN}     --ServerApp.root_dir=${JUPYTER_NOTEBOOK_DIR}     --ServerApp.allow_remote_access=1     --ServerApp.ip=0.0.0.0     --ServerApp.port=8888     --ServerApp.open_browser=False "

[Service]
TimeoutStopSec=30
ExecStartPre=/usr/bin/bash -lc 'while true; do /usr/bin/nvidia-smi -L && exit 0; sleep 1; done'

[Install]
WantedBy=default.target
